# deploy/systemd/hana-omni-connect-postgres-backup.service
[Unit]
Description=Hana Omni-Connect PostgreSQL logical backup
After=docker.service
Requires=docker.service

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/hana-omni-connect-postgres-backup
User=root
Group=root
PrivateTmp=true
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/var/backups/hana-omni-connect
